# meta-wiwynn/meta-sv300g3e/recipes-phosphor/gpio/sv300g3e-gpio-monitor-register/caterr-event.service
[Unit]
Description=IERR/MCERR Event Log Service

[Service]
Type=oneshot
ExecStart=/usr/bin/env caterr-event.sh
SyslogIdentifier=caterr-event
