# meta-wiwynn/meta-sv300g3e/recipes-phosphor/wiwynn-interrupt-dbus-register/wiwynn-interrupt-dbus-register/interrupt-dbus-register.service
[Unit]
Description=Register D-bus object for GPIO interrupt
Before=phosphor-multi-gpio-monitor.service

[Service]
Type=simple
ExecStart=/usr/bin/interruptdbusregister
SyslogIdentifier=interrupt-dbus-register
CPUSchedulingPolicy=rr
CPUSchedulingPriority=1

[Install]
WantedBy=multi-user.target

