# meta-wiwynn/meta-sv300g3e/recipes-phosphor/sv300g3e-guid/guid-generator/service/xyz.openbmc_project.guid-generator.service
[Unit]
Description=GUID Generator Service

[Service]
Type=oneshot
ExecStart=/usr/bin/guid-generator
SyslogIdentifier=guid-generator

[Install]
WantedBy=sysinit.target
