# meta-wiwynn/meta-sv300g3e/recipes-intel/asd/intel-asd/intel.asd.service
[Unit]
Description=Intel ASD Service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/asd -n eth0 -u
