# meta-wiwynn/meta-sv300g3e/recipes-phosphor/sv300g3e-powerctrl/files/host-powersoft.service
[Unit]
Description=Graceful Host Shut-down service
Requires=host-gpio.service
After=host-gpio.service

[Service]
RemainAfterExit=no
Type=oneshot
ExecStart=/usr/bin/env powersoft.sh
SyslogIdentifier=powersoft.sh
