# meta-mct/meta-f5i/recipes-phosphor/gpio/phosphor-gpio-monitor/SetAmtStart.service
[Unit]
Description=Set Advanced Memory Testing Event Starting

[Service]
Type=oneshot
ExecStart=/usr/sbin/setAmtEvent.sh assert