# meta-mct/meta-f5i/recipes-phosphor/gpio/phosphor-gpio-monitor/SetAmtFinish.service
[Unit]
Description=Set Advanced Memory Testing Event Finishing

[Service]
Type=oneshot
ExecStart=/usr/sbin/setAmtEvent.sh deassert