# meta-wiwynn/meta-sv300g3e/recipes-phosphor/sv300g3e-frb-event-handler/frb-event-handler/service/xyz.openbmc_project.frb-event-handler.service
[Unit]
Description=frb event Service
StopWhenUnneeded=false
After=interrupt-dbus-register.service

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/frb-event-handler
SyslogIdentifier=frb-event-handler

[Install]
WantedBy=multi-user.target
