# meta-wiwynn/meta-sv300g3e/recipes-phosphor/gpio/cpu-state/xyz.openbmc_project.thermtrip0.service
[Unit]
Description=assert Thermtrip0 

[Service]
Type=oneshot
ExecStart=/usr/bin/env prochot_thermtrip_update.sh assert CPU0_State thermtrip
SyslogIdentifier=thermtrip0_assert

