# meta-wiwynn/meta-sv300g3e/recipes-phosphor/sv300g3e-cpu-prochot/cpu-prochot-event/service/xyz.openbmc_project.cpu-prochot-event.service
[Unit]
Description=CPU Prochot Event Service
StopWhenUnneeded=false
After=interrupt-dbus-register.service

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/cpu-prochot-event
SyslogIdentifier=cpu-prochot-event

[Install]
WantedBy=multi-user.target
