# meta-wiwynn/meta-sv300g3e/recipes-phosphor/sv300g3e-button-handler/button-handler/service/xyz.openbmc_project.button-handler.service
[Unit]
Description=Button handler Service
StopWhenUnneeded=false
After=interrupt-dbus-register.service

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/button-handler
SyslogIdentifier=button-handler

[Install]
WantedBy=multi-user.target
