# meta-wiwynn/meta-sv300g3e/recipes-phosphor/sv300g3e-leaky-bucket/leaky-bucket/service/xyz.openbmc_project.leaky-bucket.service
[Unit]
Description=Twitter Twitter Leaky Bucket Service
StopWhenUnneeded=false

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/leaky-bucket
SyslogIdentifier=leaky-bucket

[Install]
WantedBy=multi-user.target
