# meta-wiwynn/meta-sv300g3e/recipes-phosphor/sv300g3e-service-oem/files/sv300g3e-dhcp-check.service
[Unit]
Description=Workaround for restarting the dedicated network port 
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env sv300g3e-dhcp-check.sh
SyslogIdentifier=sv300g3e-dhcp-check

[Install]
WantedBy=basic.target

