# meta-wiwynn/meta-sv300g3e/recipes-phosphor/sv300g3e-bios-updater/bios-updater/xyz.openbmc_project.bios-updater@.service
[Unit]
Description=Flash BIOS image %I to Host

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/bios-updater %i
StandardOutput=journal+console
