# meta-wiwynn/meta-sv300g3e/recipes-phosphor/sv300g3e-fast-prochot/fast-prochot-logging/service/xyz.openbmc_project.fast-prochot-logging.service
[Unit]
Description=Button handler Service
StopWhenUnneeded=false
After=interrupt-dbus-register.service

[Service]
Restart=always
RestartSec=5
ExecStart=/usr/bin/fast-prochot-logging
SyslogIdentifier=fast-prochot-logging

[Install]
WantedBy=multi-user.target
